# autoneighxy.service
[Unit]
Description=Auto IPv6 Neighbor Proxy

[Service]
ExecStart=/home/piernov/autoneighxy/autoneighxy.py

